# services/s7-public-chat.service.d/10-pinning.conf
# ═══════════════════════════════════════════════════════════════════
# s7-public-chat.service systemd drop-in — CPU pinning + memory tier
#
# This is a systemd drop-in, NOT a replacement for the main unit file.
# Drop-ins live at `<service>.service.d/*.conf` and are applied on top
# of the main unit at load time. The main `services/s7-public-chat.service`
# is NOT modified by this commit.
#
# Pattern matches s7-ollama.service (commit 3d6616c) and the 2026-04-13
# "CPU pinning + burstable like AWS T-class" brainstorm decision.
#
# Rule 3 compliance: drop-ins are non-destructive and reversible. To
# apply: `systemctl --user daemon-reload && systemctl --user restart
# s7-public-chat.service` (which is a LANDING action, not done by
# creating this file). To revert: delete this file, daemon-reload.
#
# Core assignment on i7-6600U (4 physical cores):
#   0 — s7-ollama.service        (Carli's warm model, existing pin)
#   1 — s7-bitnet-mcp.service    (future BitNet inference, post-SP1)
#   2 — s7-public-chat.service   (this drop-in — chat router + SSE streamer)
#   3 — OS + everything else     (kernel, podman, audit, breathing room)
#
# Governing rules:
#   feedback_three_rules.md   Rule 3: Protect the QUBi
#   feedback_carli_perf_followups.md  CPU pin + memory tier pattern
# ═══════════════════════════════════════════════════════════════════

[Service]
# Pin to physical core 2. Reserves the core for this service but still
# allows the kernel to burst other services onto it when idle. When
# contention happens, CPUWeight makes chat win.
AllowedCPUs=2
CPUWeight=1000

# Memory tier. Forced floor prevents the kernel from swapping out
# inference state; soft ceiling is the sustained budget; hard ceiling
# is the OOM-kill trigger. Numbers are half of Ollama's because
# public-chat does routing + SSE streaming, not inference itself.
MemoryMin=500M
MemoryHigh=2G
MemoryMax=3G

# IO priority. Ledger writes on the hot path (Rule F2 strictness: a
# row must be fsynced before the response streams). IOWeight=1000
# prioritizes chat ledger IO over everything else on core 2.
IOWeight=1000

# Explicit timeout — if the service hangs, systemd kills it rather
# than letting a stuck chat router wedge the box.
TimeoutStartSec=30
TimeoutStopSec=10
