# iac/timecapsule/s7-timecapsule-verify.service
[Unit]
Description=S7 TimeCapsule registry verification
Documentation=file:///s7/skyqubi-private/iac/timecapsule/README.md
DefaultDependencies=yes
After=local-fs.target
Before=podman.socket podman.service
ConditionPathExists=/s7/timecapsule/registry/manifest.json

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/s7-timecapsule-verify.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
